# playbooks/roles/kubernetes/server/templates/setup/apiserver.service.j2
[Service]
Environment="KUBE_API_ARGS=--advertise-address={{ ansible_host }} --allow-privileged=true --authorization-mode=Node,RBAC --client-ca-file={{ K8S_SERVER_PKI_DIR }}/{{ INT_CERT_NAME }} --enable-admission-plugins=NodeRestriction --enable-bootstrap-token-auth=true --etcd-cafile={{ K8S_SERVER_PKI_DIR }}/{{ INT_CERT_NAME }} --etcd-certfile={{ K8S_SERVER_PKI_DIR }}/{{ KUBE_APISERVER_ETCD_CLIENT_CERT_NAME }} --etcd-keyfile={{ K8S_SERVER_PKI_DIR }}/{{ KUBE_APISERVER_ETCD_CLIENT_KEY_NAME }} --etcd-servers={% for host in groups['etcd'] %}https://{{ hostvars[host].ansible_host }}:2379{{ "," if not loop.last else "" }}{% endfor %} --insecure-port=0 --kubelet-client-certificate={{ K8S_SERVER_PKI_DIR }}/{{ KUBE_APISERVER_KUBELET_CLIENT_CERT_NAME }} --kubelet-client-key={{ K8S_SERVER_PKI_DIR }}/{{ KUBE_APISERVER_KUBELET_CLIENT_CERT_NAME }} --kubelet-preferred-address-types=InternalIP,ExternalIP,Hostname --service-account-key-file={{ K8S_SERVER_PKI_DIR }}/{{ K8S_SERVICE_ACCOUNT_CERT_NAME }} --service-account-issuer=https://k8s-server.lan --api-audiences=https://k8s-server.lan --service-account-signing-key-file={{ K8S_SERVER_PKI_DIR }}/{{ K8S_SERVICE_ACCOUNT_KEY_NAME }} --service-cluster-ip-range=10.96.0.0/12 --tls-cert-file={{ K8S_SERVER_PKI_DIR }}/{{ KUBE_APISERVER_CERT_NAME }} --tls-private-key-file={{ K8S_SERVER_PKI_DIR }}/{{ KUBE_APISERVER_KEY_NAME }}"
ExecStart=
ExecStart=/usr/bin/kube-apiserver $KUBE_API_ARGS
